# platforms/raspberry-pi2/data/dartino-agent
#!/bin/bash
# Start/stop the Dartino agent.
#
### BEGIN INIT INFO
# Provides:          dartino-agent
# Required-Start:    $syslog $network
# Required-Stop:     $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Dartino agent for managing Dartino VMs.
# Description:       The Dartino agent is used by the Dartino compiler to start
#                    and stop Dartino VMs on a Raspberry Pi2 device.
### END INIT INFO

PATH=/bin:/usr/bin:/sbin:/usr/sbin
AGENT_SNAPSHOT=/usr/lib/dartino-agent/bin/fletch-agent.snapshot
NAME=dartino-agent

if [ ! -r /lib/lsb/init-functions ]; then
  echo "Missing lsb init-functions to run the Fletch agent"
  exit 1
fi
. /lib/lsb/init-functions

# Setup environment variables until we support passing arguments to fletch
# programs.
if [ ! -r /etc/default/dartino-agent ]; then
  echo "Missing Dartino environment file, /etc/default/dartino-agent"
  exit 1
fi
. /etc/default/dartino-agent

# Check the environment and create VM directories if necessary.
if [ ! -x $FLETCH_VM ]; then
  log_daemon_msg "No Dartino VM found or VM not executable"
fi
if [ ! -d $VM_LOG_DIR ]; then
  mkdir -p $VM_LOG_DIR
fi

case "$1" in
start)	log_daemon_msg "Starting" "$NAME"
	start-stop-daemon --background --start --quiet --oknodo \
	    --pidfile $AGENT_PID_FILE --exec $FLETCH_VM $AGENT_SNAPSHOT
	log_end_msg $?
	;;
stop)	log_daemon_msg "Stopping" "$NAME"
	killproc -p $AGENT_PID_FILE $FLETCH_VM
	RETVAL=$?
	[ $RETVAL -eq 0 ] && [ -e "$AGENT_PID_FILE" ] && rm -f $AGENT_PID_FILE
	log_end_msg $RETVAL
	;;
restart) log_daemon_msg "Restarting" "$NAME"
	$0 stop
	sleep 2
	$0 start
	;;
status)
	status_of_proc -p $AGENT_PID_FILE $FLETCH_VM $NAME && exit 0 || exit $?
	;;
*)	log_action_msg \
	    "Usage: /etc/init.d/dartino-agent {start|stop|status|restart}"
	exit 2
	;;
esac
exit 0
